# root/etc/profile.d/dircolors.sh
#!/bin/bash

eval `dircolors /etc/dircolors/dircolors.ansi-dark`
